# debian/control/postinst
#!/bin/sh
# postinst script for #PACKAGE#
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


create_user() {
  if ! id odn-restsql > /dev/null 2>&1 ; then
    adduser --system --no-create-home --disabled-password --home /usr/share/odn-restsql odn-restsql
  fi
}


service_start() {
    if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d odn-restsql start || true
    else
        service odn-restsql start || true
    fi
}

prepare_dir() { 
    mkdir -p /usr/share/odn-restsql/webapps
    mkdir -p /var/share/odn-restsql
    mkdir -p /var/log/odn-restsql
    mkdir -p /var/odn-simple/restsql
    mkdir -p /var/lib/odn-restsql
    mkdir -p /var/cache/odn-restsql
    
    
    chown -R odn-restsql /usr/share/odn-restsql/webapps /var/share/odn-restsql /var/log/odn-restsql /var/odn-simple/restsql /var/lib/odn-restsql /var/cache/odn-restsql  /etc/odn-simple/restsql/
    ln -s /etc/odn-restsql/tomcat7 /usr/share/odn-restsql/conf
    ln -s /var/log/odn-restsql/ /usr/share/odn-restsql/logs
    ln -s /var/tmp/ /usr/share/odn-restsql/temp
    ln -s /var/cache/odn-restsql  /usr/share/odn-restsql/work
    
    chmod +x  /etc/init.d/odn-restsql
}

case "$1" in
    configure)
    create_user
    prepare_dir
    service_start
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

